# nodejs/etc/systemd/system/nodejs-watcher.service
[Unit]
Description=MyApp Server Watcher

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart app

[Install]
WantedBy=multi-user.target