# nodejs/etc/systemd/system/nodejs-watcher.path
[Path]
PathModified=/home/node/app

[Install]
WantedBy=multi-user.target